# Gateway/taskplanner.service
[Unit]
Description=TaskPlanner
After=network.target

[Service]
ExecStart=/usr/bin/node ~/TaskPlanner/Gateway/taskplanner.js

[Install]
WantedBy=multi-user.target
